# debian/cnb-cache.timer
[Unit]
Description=Run cnb-cache-update every midnight

[Timer]
OnCalendar=00:01
Persistent=true

[Install]
WantedBy=timers.target
